# tamperedchef_pdfeditor_run_persistence.yml
title: TamperedChef PDFEditorUpdater Run Key Persistence
id: 609d7232-b537-4fc7-b645-6b778e1bea3a
status: experimental
description: Detects the PDFEditorUpdater Run value and --cm=--fullupdate command
  used by the AppSuite PDF Editor TamperedChef cluster.
references:
- https://blog.lukeacha.com/2025/06/suspicious-recipe-app.html
- https://blog.lukeacha.com/2025/08/pdf-editor-analysis.html
author: Luke Acha
date: '2026-07-22'
tags:
- attack.persistence
- attack.t1060
- attack.t1547.001
logsource:
  category: registry_set
  product: windows
detection:
  selection_key:
    TargetObject|endswith: \Software\Microsoft\Windows\CurrentVersion\Run\PDFEditorUpdater
  selection_value:
    Details|contains|all:
    - PDF Editor
    - --cm=--fullupdate
  condition: selection_key and selection_value
falsepositives:
- Unlikely
level: high
